# contrib/systemd/rebuilderd-stats.timer
[Unit]
Description=rebuilderctl stats: collect a daily snapshot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
