# contrib/systemd/rebuilderd-stats.service
[Unit]
Description=rebuilderctl stats: collect a daily snapshot

[Service]
User=rebuilderd
ExecStart=/usr/bin/rebuildctl stats collect
