# contrib/systemd/rebuilderd-worker@.service
[Unit]
Description=rebuilderd-worker: rebuild packages

[Service]
ExecStart=/usr/bin/rebuilderd-worker -n %I connect

[Install]
WantedBy=multi-user.target
