# contrib/systemd/rebuilderd-sync@.timer
[Unit]
Description=rebuilderctl sync: periodically import packages

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
