# contrib/systemd/rebuilderd-sync@.service
[Unit]
Description=rebuilderctl sync: periodically import packages

[Service]
User=rebuilderd
ExecStart=/usr/bin/rebuildctl pkgs sync-profile %I
